# deploy/maotai-api.service
[Unit]
Description=Maotai Reserve Admin API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/maotai
Environment=PYTHONPATH=/opt/maotai
EnvironmentFile=-/opt/maotai/deploy/.env
ExecStart=/opt/maotai/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --app-dir /opt/maotai/backend
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
